# deploy/systemd/mcp-hugo-server-go.service
# Unified mode: supports content.write, site.admin, system.admin tools.
# ReadWritePaths must include site_root, content_root, and the OAuth storage
# path (if using json/sqlite backend). Adjust to your deployment paths.
#
# To run in read-only mode (anonymous + content.read only): remove the
# ReadWritePaths lines and re-add ProtectSystem=strict.

[Unit]
Description=Hugo unified MCP server
Documentation=https://github.com/jmrGrav/mcp-hugo-server-go
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mcp-hugo-server-go
Group=mcp-hugo-server-go
Environment=MCP_HUGO_SERVER_CONFIG=/etc/mcp-hugo-server-go/config.yaml
ExecStart=/usr/local/bin/mcp-hugo-server-go
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
# ProtectSystem=full allows writes to /var and /run but not /etc or /usr.
# Use "strict" here only for read-only mode; write/admin mode needs write access.
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Grant write access to the Hugo site directories and OAuth token store.
# Adjust these paths to match your deployment (site_root, content_root,
# and oauth.storage_path from config.yaml).
ReadWritePaths=/var/lib/mcp-hugo-server-go /srv/hugo-site
# Hugo needs to read layouts, static, and themes too.
ReadOnlyPaths=/etc/mcp-hugo-server-go

[Install]
WantedBy=multi-user.target
