# deploy/systemd/mcp-hugo-server-go.service
[Unit]
Description=Hugo unified MCP server
Documentation=https://github.com/jmrGrav/mcp-hugo-server-go
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mcp-hugo-server-go
Group=mcp-hugo-server-go
Environment=MCP_HUGO_SERVER_CONFIG=/etc/mcp-hugo-server-go/config.yaml
ExecStart=/usr/local/bin/mcp-hugo-server-go
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
# ReadOnlyPaths and ReadWritePaths are site-specific.
# Set them directly in this file or via an override.conf drop-in.
# build_site and preview_build require write access to hugo_root/public
# and hugo_root/resources in addition to content_root.
#
# Example (inline or override.conf):
#   ReadOnlyPaths=/etc/mcp-hugo-server-go
#   ReadWritePaths=/var/lib/mcp-hugo-server-go \
#       /your/hugo-site/content \
#       /your/hugo-site/resources \
#       /your/hugo-site/public
#   Environment=PATH=/usr/local/bin:/usr/bin:/bin
#
# Also ensure the service user has group-write permission on those directories:
#   chown -R user:servicegroup /your/hugo-site/public /your/hugo-site/resources
#   chmod -R ug+w /your/hugo-site/public /your/hugo-site/resources

[Install]
WantedBy=multi-user.target
